# log-viewer/rootfs/etc/services.d/logviewer/run
#!/usr/bin/with-contenv bashio
# ==============================================================================
# Community Hass.io Add-ons: Log Viewer
# Runs the Log Viewer
# ==============================================================================
readonly LOG_FILE="/config/home-assistant.log"
declare -a options

bashio::log.info "Starting Log Viewer..."

options+=(--user "$(bashio::config 'username')")
options+=(--password "$(bashio::config 'password')")

if bashio::config.true 'ssl'; then
    options+=(--ssl)
    options+=(--certfile "/ssl/$(bashio::config 'certfile')")
    options+=(--keyfile "/ssl/$(bashio::config 'keyfile')")
fi

bashio::log.debug "$(bashio::config 'filters')"

exec ws-log "${options[@]}" "${LOG_FILE}"
