# log-viewer/rootfs/etc/services.d/logviewer/run
#!/usr/bin/with-contenv bash
# ==============================================================================
# Community Hass.io Add-ons: Log Viewer
# Runs the Log Viewer
# ==============================================================================
# shellcheck disable=SC1091
source /usr/lib/hassio-addons/base.sh

readonly LOG_FILE="/config/home-assistant.log"
declare -a options

options+=(--user "$(hass.config.get 'username')")
options+=(--password "$(hass.config.get 'password')")

if hass.config.true 'ssl'; then
    options+=(--ssl)
    options+=(--certfile "/ssl/$(hass.config.get 'certfile')")
    options+=(--keyfile "/ssl/$(hass.config.get 'keyfile')")
fi

hass.log.debug "$(hass.config.get 'filters')"

exec ws-log "${options[@]}" "${LOG_FILE}"
