# log-viewer/rootfs/etc/services.d/logviewer/run
#!/usr/bin/with-contenv bashio
# ==============================================================================
# Home Assistant Community Add-on: Log Viewer
# Runs the Log Viewer
# ==============================================================================
bashio::log.info "Starting Log Viewer..."
exec ws-log "/config/home-assistant.log"
